# ubuntu-base/files/autoresize.service
[Unit]
Description=Resize disk on boot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/autoresize

[Install]
WantedBy=default.target
